# com.trii.slack-listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trii.slack-listener</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/georgeclark/Desktop/Projects.nosync/TRII/.venv/bin/python3</string>
        <string>/Users/georgeclark/Desktop/Projects.nosync/TRII/slack-listener.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/georgeclark/Desktop/Projects.nosync/TRII</string>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/georgeclark/Desktop/Projects.nosync/TRII/session-log/slack-listener-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/georgeclark/Desktop/Projects.nosync/TRII/session-log/slack-listener-stderr.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
